# ansible/roles/autorestic/templates/autorestic-cloud.timer.j2
[Unit]
Description=Run Autorestic cloud sync (3 AM)

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
